# installer/linux/ionconnect-core.service
[Unit]
Description=IonConnect core (mouse/keyboard/clipboard sharing service)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Ajustar la ruta si `ionconnect-core` no quedó en /usr/bin (p. ej.
# ~/.local/bin/ionconnect-core si se instaló con install.sh).
ExecStart=%h/.local/bin/ionconnect-core
Restart=on-failure
RestartSec=2
# Servicio de usuario, no de sistema: corre con la sesión gráfica del
# usuario (necesita acceso a X11/Wayland, XDG portals, D-Bus de sesión).
# Instalar en ~/.config/systemd/user/, habilitar con:
#   systemctl --user enable --now ionconnect-core.service

[Install]
WantedBy=default.target
